# avatar-tank-advanced.service
[Unit]
Description=Avatar Tank Advanced System - Complete MediaMTX Integration
After=network.target
Wants=network.target

[Service]
Type=simple
User=havatar
Group=havatar
WorkingDirectory=/home/havatar/Avatar-robot
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 modules/mediamtx_main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/havatar/Avatar-robot
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1

# Security settings (minimal for development)
NoNewPrivileges=false
PrivateTmp=false

# Allow access to devices
SupplementaryGroups=video audio dialout

[Install]
WantedBy=multi-user.target


