# avatar-tank.service
[Unit]
Description=Avatar Tank System
After=network-online.target sound.target
Wants=network-online.target sound.target

[Service]
Type=simple
User=havatar
Group=havatar
WorkingDirectory=/home/havatar/Avatar-robot
ExecStart=/bin/bash /home/havatar/Avatar-robot/start_avatar_simple.sh
Restart=always
# Exponential-style backoff via StartLimit* + growing delays in script; base here
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=mixed
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
WatchdogSec=0
NotifyAccess=none
# Resource limits for reliability on 4G devices
MemoryMax=350M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
